# etc/systemd/user/bing-dl.service
[Unit]
Description=bing-dl Daemon
Wants=bing-dl.timer

[Service]
Type=oneshot
ExecStart=%h/.local/bin/bing-dl
