# etc/systemd/user/bing-dl.service
[Unit]
Description=bing-dl Daemon

[Service]
Type=oneshot
ExecStart=%h/.local/bin/bing-dl

[Install]
WantedBy=default.target
